# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:xinit]
command=xinit /etc/X11/xinit/Xsession openbox -- /usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ardelay 200 -arinterval 30
autorestart=true
user=alpine
priority=100

[program:x11vnc]
command=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :1 -repeat -nopw -wait 5 -permitfiletransfer -tightfilexfer
user=alpine
autorestart=true
priority=200
